# scripts/services/hero2-demo/hero2-battery-manager-laptop.service
[Unit]
Description=running battery_manager_laptop

[Install]
WantedBy=multi-user.target

[Service]
User=demo
Restart=on-failure
ExecStart=/bin/bash -c 'source /home/demo/.tue/setup.bash && \
                        roslaunch hero_bringup battery_manager_laptop
