# scripts/hero1-display.service
[Unit]
Description=Head display on Hero's head
After=graphical.target
Wants=graphical.target

[Install]
WantedBy=multi-user.target

[Service]
User=hsr-hmi
Restart=on-abort
Environment=SYSTEMD_LOG_LEVEL=debug
Environment='DISPLAY=:0'
ExecStart=/usr/bin/hero-display.AppImage
