# scripts/ntpdate.service
[Unit]
Description=NTPDATE

[Install]
WantedBy=hero-roscore.service

[Service]
Restart=on-abort
ExecStart=/usr/sbin/ntpdate ntp.ubuntu.com
